# casabot/rootfs/etc/services.d/casabot/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# CasaBot finish script
# ==============================================================================

bashio::log.info "CasaBot stopped"