# casabot/rootfs/etc/s6-overlay/s6-rc.d/init-casabot/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: CasaBot
# Configures CasaBot before running
# ==============================================================================
declare port

# Migrate add-on data from the Home Assistant config folder,
# to the add-on configuration folder.



