# casabot/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Nginx finish script
# ==============================================================================

bashio::log.info "Nginx stopped"